# upstart/dinopix.conf
description "dinopix (dinosaurpictures.org)"
start on filesystem and net-device-up IFACE=eth0
stop on shutdown

exec sudo -u ian "/home/ian/dinopix/upstart/upstart.sh" >>/var/log/upstart/dinopix.log
